# config/ethbot.service
[Unit]
Description=ETH Strategy Dashboard - FastAPI Service
After=network.target
Documentation=https://github.com/yourusername/ethbot

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/var/www/ethbot

# Environment
Environment="PATH=/var/www/ethbot/venv/bin"
EnvironmentFile=/var/www/ethbot/.env

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/ethbot

# Start command
ExecStart=/var/www/ethbot/venv/bin/gunicorn \
    -w 4 \
    -k uvicorn.workers.UvicornWorker \
    --bind 127.0.0.1:8000 \
    --timeout 60 \
    --max-requests 1000 \
    --max-requests-jitter 50 \
    --access-logfile /var/log/ethbot/access.log \
    --error-logfile /var/log/ethbot/error.log \
    --log-level info \
    app:app

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
